# systemd/saleapp-scheduler.service.d/10-xray.conf
# Drop-in override: couple the scheduler to the xray bridge.
#
# Installed as /etc/systemd/system/saleapp-scheduler.service.d/10-xray.conf.
# Intentionally overlays the existing scheduler unit (which wraps python3 in
# xvfb-run for Chromium-based cart-add) instead of replacing it, so the
# existing ExecStart, xvfb setup, and env are preserved verbatim.
#
# The only change is to declare saleapp-xray.service as a hard dependency,
# so the bridge comes up before the scheduler and stops after it.
[Unit]
After=network-online.target saleapp-xray.service
Wants=network-online.target
Requires=saleapp-xray.service
