# systemd/saleapp-scheduler.service.d/10-xray.conf
# Drop-in override: couple the scheduler to the xray bridge.
#
# Installed as /etc/systemd/system/saleapp-scheduler.service.d/10-xray.conf.
# Intentionally overlays the existing scheduler unit (which wraps python3 in
# xvfb-run for Chromium-based cart-add) instead of replacing it, so the
# existing ExecStart, xvfb setup, and env are preserved verbatim.
#
# v1.26 Phase 84.5: changed `Requires=` to `Wants=` because every
# `systemctl reload-or-restart saleapp-xray` (issued by
# vless.manager._reload_xray_systemd whenever the pool admits new nodes)
# was cascading into a scheduler restart. With Requires=, systemd
# stop+start the scheduler around every xray reload, which killed mid-run
# scrapes and added ~90s of downtime per pool refresh — directly
# contributing to the user-reported "Обновлено: 19 мин" staleness.
# Wants= keeps the "xray comes up before scheduler" ordering (via
# `After=`) without the restart cascade. The scheduler is robust to
# transient xray downtime — it just retries when xray is back.
[Unit]
After=network-online.target saleapp-xray.service
Wants=network-online.target saleapp-xray.service
