# systemd/saleapp-backend.service.d/10-xray.conf
[Unit]
After=network-online.target saleapp-xray.service
Wants=network-online.target
Requires=saleapp-xray.service
