# systemd/saleapp-backend.service.d/10-xray.conf
# Drop-in override: couple the backend to the xray bridge.
#
# Installed as /etc/systemd/system/saleapp-backend.service.d/10-xray.conf.
# Intentionally overlays the existing backend unit so the existing
# ExecStart and env are preserved verbatim.
#
# v1.26 Phase 84.5: changed `Requires=` to `Wants=` (see matching note in
# saleapp-scheduler.service.d/10-xray.conf). xray reload-or-restart
# events were cascading into backend restarts, breaking SSE connections
# and stalling /api/products responses for ~90s per pool refresh. The
# backend is robust to xray downtime via the existing breaker /
# fallback paths.
[Unit]
After=network-online.target saleapp-xray.service
Wants=network-online.target saleapp-xray.service
